# config/newreleases/newreleases.service
[Unit]
Description=Git releases Watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/newreleases
