# newreleases/newreleases.timer
[Unit]
Description=Git releases Watchdog

[Timer]
OnCalendar=7..23:0/15
Persistent=true

[Install]
WantedBy=timers.target
