# newreleases/newreleases.timer
[Unit]
Description=Git releases Watchdog

[Timer]
OnCalendar=7..22/3:0
Persistent=true

[Install]
WantedBy=timers.target
