# update/update.service
[Unit]
Description=Update all package bases

[Service]
ExecStart=/usr/bin/updpkgs

[Install]
WantedBy=default.target

