# wireguard/stat/wgstat.timer
[Unit]
Description=Timer for WgStat Logger

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
