# wireguard/stat/wgstat.service
[Unit]
Description=Dump transfer log for wg0 interface each hour 

[Service]
ExecStart=/usr/bin/wg-stat

[Install]
WantedBy=default.target

